# 99-ddcci-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/ddcci_probe.sh"
